# packages/ceti-tag-data-capture/debian/ceti-tag-data-capture.service
[Unit]
Description=Captures data from input sensors of the Whale Tag and saves it onto the filesystem
After=ceti-tag-set-hostname.service
StartLimitIntervalSec=500
StartLimitBurst=5

[Service]
Type=exec
ExecStart=/opt/ceti-tag-data-capture/bin/cetiTagApp
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
