# packages/ceti-tag-data-capture/debian/ceti-tag-data-capture.service
[Unit]
Description=Captures data from input sensors of the Whale Tag and saves it onto the filesystem
After=network.target local-fs.target getty.target

[Service]
Type=simple
WorkingDirectory=/opt/ceti-tag-data-capture/ipc
ExecStart=/opt/ceti-tag-data-capture/ipc/tagMonitor.sh
Restart=always
RestartSec=60

[Install]
WantedBy=multi-user.target
