# packages/ceti-tag-data-capture/debian/ceti-tag-data-capture.service
[Unit]
Description=Captures data from input sensors of the Whale Tag and saves it onto the filesystem

[Service]
Type=simple
WorkingDirectory=/opt/ceti-tag-data-capture/ipc
ExecStart=/opt/ceti-tag-data-capture/ipc/cetiTagMonitor.sh
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target