# packages/ceti-tag-set-hostname/debian/ceti-tag-set-hostname.service
[Unit]
Description=Set hostname of the system to UID of the device

[Service]
Type=oneshot
ExecStart=-/opt/ceti-tag-set-hostname/ceti-tag-set-hostname.sh
RemainAfterExit=no

[Install]
WantedBy=multi-user.target
