# packages/ceti-tag-battery-monitor/debian/ceti-tag-battery-monitor.service
[Unit]
Description=Monitors the remaining charge of the battery and shuts down the system if the battery is low
After=network.target local-fs.target getty.target

[Service]
Type=simple
WorkingDirectory=/opt/ceti-tag-battery-monitor/bin
ExecStart=/opt/ceti-tag-battery-monitor/bin/batterymon.sh
Restart=always
RestartSec=60

[Install]
WantedBy=multi-user.target
